# tutorial/templates/supervisor-fastapi.conf
; Supervisor config untuk FastAPI (uvicorn)
; Copy ke: /etc/supervisor/conf.d/hunter-backend.conf
; Lalu: sudo supervisorctl reread && sudo supervisorctl update

[program:hunter-backend]
command=/var/www/hunter.elexart.com/backend/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8001 --workers 2
directory=/var/www/hunter.elexart.com/backend
user=deploy
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/hunter-backend.out.log
stderr_logfile=/var/log/supervisor/hunter-backend.err.log
environment=
    PATH="/var/www/hunter.elexart.com/backend/venv/bin:/usr/bin",
    PYTHONUNBUFFERED="1"
